# option-bridge/claude-lark-refresh.timer
[Unit]
Description=Claude Code - lark-cli token refresh timer (every 5 days)
After=network-online.target

[Timer]
# Every 5 days at 10:00 AM local time
OnCalendar=*-*-01,06,11,16,21,26 10:00:00
Persistent=true
RandomizedDelaySec=1800

[Install]
WantedBy=timers.target
